# apt-ack.service
[Unit]
Description=APT-ACK Threat Intelligence Collection
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=%u
WorkingDirectory=/home/baals/Desktop/bigdev/APT-ACK
Environment="PATH=/home/baals/Desktop/bigdev/APT-ACK/venv/bin:/usr/local/bin:/usr/bin"
ExecStart=/home/baals/Desktop/bigdev/APT-ACK/venv/bin/python /home/baals/Desktop/bigdev/APT-ACK/scripts/run_collection.py --skip-nvd
StandardOutput=append:/home/baals/Desktop/bigdev/APT-ACK/logs/collection.log
StandardError=append:/home/baals/Desktop/bigdev/APT-ACK/logs/collection.log

[Install]
WantedBy=multi-user.target